# ev3dev-nap_0.1-1_all/DEBIAN/preinst
#!/bin/bash
set -e

# stop django server
if [ -f /etc/init.d/nap ]
then
  invoke-rc.d ev3dev-nap stop
fi