# ev3dev-nap_0.1-1_all/DEBIAN/prerm
#!/bin/bash
set -e

# stop ev3dev-nap server
invoke-rc.d ev3dev-nap stop
ifdown tether